# root/etc/services.d/influxdb/run
#!/usr/bin/with-contenv bash

s6-setuidgid influxdb exec influxd -config /etc/influxdb.conf
